# supervisord.conf
[supervisord]
nodaemon=true

[program:avahi]
command=/usr/sbin/avahi-daemon --no-drop-root

[program:smb]
command=/usr/local/samba/sbin/smbd -F -d 1 -s /etc/samba/smb.conf
